# com.nologik.awdl-monitor.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.nologik.awdl-monitor</string>
    
    <key>ProgramArguments</key>
    <array>
        <string>/Users/nologik/awdl-mac-monitor/monitor_awdl_live.sh</string>
        <string>30</string>
    </array>
    
    <key>RunAtLoad</key>
    <true/>
    
    <key>KeepAlive</key>
    <true/>
    
    <key>StandardOutPath</key>
    <string>/Users/nologik/awdl-mac-monitor/monitor.log</string>
    
    <key>StandardErrorPath</key>
    <string>/Users/nologik/awdl-mac-monitor/monitor.error.log</string>
    
    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
    </dict>
    
    <key>ProcessType</key>
    <string>Background</string>
    
    <key>ThrottleInterval</key>
    <integer>30</integer>
    
    <key>Nice</key>
    <integer>1</integer>
</dict>
</plist>
